# launchd/com.ryuryu.homepage-engine-daily-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ryuryu.homepage-engine-daily-report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/okamotoryuya/.claude/bots/lib/daily-report-core.sh</string>
        <string>/Users/okamotoryuya/.claude/bots/homepage-engine/config.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/okamotoryuya/repository/homepage-engine/logs/daily-report-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/okamotoryuya/repository/homepage-engine/logs/daily-report-stderr.log</string>
    <key>RunAtLoad</key>
    <false/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/okamotoryuya/.local/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/okamotoryuya</string>
        <key>TZ</key>
        <string>Asia/Tokyo</string>
    </dict>
</dict>
</plist>
